# s6-config/prod/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bash

unset S6_KEEP_ENV
unset PHP_TIMEZONE
unset PHP_UPLOAD_MAX_FILESIZE
unset PHP_POST_MAX_SIZE
unset PHP_MEMORY_LIMIT
unset PHP_MAX_EXECUTION_TIME
unset PHP_OPCACHE_ENABLED
unset PHP_OPCACHE_MEMORY_CONSUMPTION
unset PHP_OPCACHE_JIT_BUFFER_SIZE
unset PHP_OPCACHE_JIT
unset PHP_OPCACHE_VALIDATE_TIMESTAMPS
unset PHP_OPCACHE_REVALIDATE_FREQ

echo "Starting php-fpm"

exec php-fpm -F