# s6-config/dev/etc/s6-overlay/s6-rc.d/fix-permissions/up
/etc/s6-overlay/s6-rc.d/fix-permissions/fix-permissions.sh