# s6-config/prod/etc/s6-overlay/s6-rc.d/php-ini-override/up
gomplate --input-dir=/usr/local/etc/php/conf.d.templates --output-dir=/usr/local/etc/php/conf.d